# templates/dashboard.service.j2
[Unit]
After=network.service

[Service]
WorkingDirectory={{ dashboard_dir }}
ExecStart={{ dashboard_dir }}/dashboard
Restart=always

[Install]
WantedBy=default.target
